# test/ansible/templates/systemd/agate.service.j2
# 2018-12-23 (cc) <paul4hough@gmail.com> -*- conf-unix -*-
#
[Unit]
Description=prometheus alertmanager alert gateway
Documentation=https://github.com/pahoughton/agate
After=network-online.target

[Service]
User={{ remed_user }}
ExecStart={{ remed_dir }}/bin/agate \
	     --addr ":{{ port_agate }}" \
	     --config {{ remed_dir }}/config/agate.yml \
	     --data {{ remed_dir }}/data \
	     --debug

[Install]
WantedBy=multi-user.target
